# generators/app/templates/alpine-nginx-theia-ide-hl/etc/services.d/nginx/run
#!/command/with-contentenv bash

exec /usr/sbin/nginx -c /config/nginx/nginx.conf

