# generators/app/templates/alpine-binary-c-hl/etc/services.d/unitd/run
#!/usr/bin/with-contenv sh

echo "starting unit daemon (unit.nginx)"
exec /usr/sbin/unitd --no-daemon --log /dev/stderr
