# generators/app/templates/alipne-nginx-theia-ide/etc/services.d/php-fpm/run
#!/bin/bash

# exec /usr/sbin/php-fpm7 -F
source /etc/hluser
s6-setuidgid $HL_USER_USERNAME /usr/sbin/php-fpm7 -F
