# generators/app/templates/alpine-binary-c-dynamic-hl/etc/services.d/unitd/run
#!/command/with-contenv bash

echo "starting unit daemon (unit.nginx)"
exec /usr/sbin/unitd --no-daemon --log /dev/stderr

